# ansible/roles/repos/files/certbot-renewal.service
[Unit]
Description=Certbot Renewal

[Service]
ExecStart=/usr/bin/certbot renew --webroot --webroot-path /repos/ --post-hook "systemctl reload httpd"
